# crontab
*/1 * * * * /bin/bash /root/ping.sh >> lastrun.log
